# docker/service/link-tagging/run
#!/bin/bash

/usr/bin/java -cp /var/www/app/coop-all.jar com.netegreek.chattr.CoopApplication server /var/www/app/config.yml
